# agent/dist/linux-x64/netflowAgent.service
[Unit]
Description=netflowAgent NetFlow export
Documentation=https://github.com/kolixxx/nflowCapAgent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/netflowAgent/netflowAgent --config /etc/netflowAgent/config.toml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
